# ansible/roles/minio-exports/templates/bezaforge-minio-exports.timer.j2
[Unit]
Description=Schedule for BezaForge nightly MinIO bucket exports
Documentation=https://github.com/thejollydev/bezaforge-infrastructure

[Timer]
OnCalendar={{ minio_export_schedule }}
# If we missed a run (host was off), fire immediately on boot.
Persistent=true
# Small jitter prevents lockstep firing if many timers share the schedule.
RandomizedDelaySec=2min
Unit=bezaforge-minio-exports.service

[Install]
WantedBy=timers.target
