# ansible/roles/forge-brizza/templates/hermes-backup-rsync.service.j2
# {{ ansible_managed }}
[Unit]
Description=Nightly rsync of ~/.hermes → bezapool/brizza-backup (FORGE-25)
# Order after the insurance exports so a Persistent catch-up after
# downtime runs export-then-backup in sequence, same as a normal night.
After=hermes-insurance.service network-online.target
RequiresMountsFor={{ forge_brizza_hermes_backup_mount }}

[Service]
Type=oneshot
# Root: ~/.hermes contents are uid-10000; the NFS export is
# no_root_squash so ownership is preserved raw.
ExecStart=/usr/local/bin/hermes-backup-rsync.sh
StandardOutput=journal
StandardError=journal
# Generous cold-run ceiling. With caches excluded the steady-state
# payload is small and runs in seconds; the margin only matters for a
# first cold copy over the sync NFS export.
TimeoutStartSec=2h
