# ansible/roles/forge-brizza/templates/hermes-insurance.service.j2
# {{ ansible_managed }}
[Unit]
Description=Brizza portability exports (sessions JSONL + quick state snapshot — BRIZZA-58)
# Bare-metal since BRIZZA-66 (BRIZZA-73): no docker dependency. The export
# reads ~/.hermes directly and `hermes backup -q` snapshots live state
# consistently, so it needs no ordering against docker or brizza.service.

[Service]
Type=oneshot
User={{ forge_brizza_docker_user }}
ExecStart=/usr/local/bin/hermes-insurance.sh
