# ansible/roles/forge-brizza/templates/hermes-backup-rsync.service.j2
# {{ ansible_managed }}
[Unit]
Description=Nightly rsync of ~/.hermes → bezapool/brizza-backup over SSH (FORGE-25)
# Order after the insurance exports so a Persistent catch-up after
# downtime runs export-then-backup in sequence, same as a normal night.
After=hermes-insurance.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Root: ~/.hermes contents are uid-10000, and the push uses root's
# dedicated SSH key to reach root on the hypervisor.
ExecStart=/usr/local/bin/hermes-backup-rsync.sh
StandardOutput=journal
StandardError=journal
# Small payload over SSH — runs in seconds. Generous ceiling for a cold
# copy only.
TimeoutStartSec=1h
